# source/plgcheck/etc/cron.daily/plgcheck.sh
#!/bin/bash
/usr/local/emhttp/plugins/plgcheck/plgcheck.php

